# streaming-sync.service
[Unit]
Description=TMDB Playlist Background Sync Daemon
After=network.target nginx.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/streaming
ExecStart=/usr/bin/php /var/www/streaming/background_sync_daemon.php --daemon
Restart=always
RestartSec=30
StandardOutput=append:/var/www/streaming/logs/sync_daemon.log
StandardError=append:/var/www/streaming/logs/sync_daemon.log

[Install]
WantedBy=multi-user.target
